# publishes/docker/php/supervisor.dev.ini
[supervisord]
nodaemon=true
user=root
pidfile=/var/run/supervisord.pid

[program:php-fpm]
process_name=%(program_name)s
command=php-fpm
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:listen --sleep=3 --tries=3
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=8
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopwaitsecs=3600

[program:crond]
process_name=%(program_name)s
command=crond -f
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
